# build/step6/copy/etc/systemd/system/run-autostart.service
[Unit]
Description=run-autostart
After=basic.target
After=network.target
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=true
WorkingDirectory=/root
ExecStart=/root/run-autostart.sh
KillMode=process

[Install]
WantedBy=multi-user.target
